# .config/systemd/user/rclone-upload-spotify-history.timer
[Unit]
Description=Upload Spotify local history file every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
