# .config/systemd/user/restic-backup-home.timer
[Unit]
Description=Start restic backup

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
